# HackerOS-Events.service
[Unit]
Description=HackerOS Events

[Service]
Type=oneshot
ExecStart=/usr/lib/HackerOS/HackerOS-Events
Restart=no

[Install]
WantedBy=default.target
